# deploy/vacancy-fetch.timer
[Unit]
Description=Vacancy Aggregator: обход по будням в 7:00 МСК

[Timer]
# Пояс ЯВНО. Сервер живёт в UTC: строчка "07:00" без пояса дала бы 10:00 по
# Москве, и утренний дайджест приходил бы к обеду. У Москвы нет перехода на
# летнее время с 2014-го, но полагаться на это незачем — пусть systemd считает.
#
# Mon-Fri: в выходные HR не читают, а понедельничный прогон накроет субботу с
# воскресеньем — LIMIT=50 на канал против потока ~4 поста в сутки.
OnCalendar=Mon-Fri 07:00 Europe/Moscow

# Сервер лежал в 7:00 и поднялся в 9:00 — таймер догонит пропущенный запуск.
# Крон в этом случае просто промолчал бы, и ты не получил бы дайджест, не
# понимая почему.
Persistent=true

# Чтобы девять каналов не получали запрос секунда в секунду каждый день
RandomizedDelaySec=120

[Install]
WantedBy=timers.target
